# templates/home-assistant.service.j2
[Unit]
Description=Home Assistant
After=network.target

[Service]
Type=simple
User={{ ha_user }}
ExecStart={{ hass_command.stdout }} --config {{ ha_conf_dir }}

[Install]
WantedBy=multi-user.target
