# templates/home-assistant.service.j2
[Unit]
Description=Home Assistant
After=network.target

[Service]
Type=simple
User=%i
ExecStart=/usr/bin/hass --config {{ ha_conf_dir }}

[Install]
WantedBy=multi-user.target

